# crontab
0 6 * * * /react/argovis/generate_argo_sitemap.sh
